# src/main/config/hive-auto-partitioner.ini
[program:hive-auto-partitioner]
command = /opt/daplab/hive-auto-partitioner/bin/run.sh --topic trumpet --zk.connect localhost:2181 --configFile /opt/daplab/hive-auto-partitioner/config/hive-partitions.json
autostart = true
autorestart = true
startsecs = 10
startretries = 3
user = hdfs
stopsignal = KILL
stopasgroup = true
stdout_logfile = /opt/daplab/hive-auto-partitioner/logs/supervisor.stdout.log
stdout_logfile_maxbytes = 100000000
stdout_logfile_backups = 5
stderr_logfile = /opt/daplab/hive-auto-partitioner/logs/supervisor.stderr.log
stderr_logfile_maxbytes = 100000000
stderr_logfile_backups = 5
